# assets/etc/supervisor/conf.d/apache.ini
[program:apache2]
command=/bin/bash -c "mkdir -p /var/log/apache2 ; \
                      touch /var/log/apache2/access.log; \
                      touch /var/log/apache2/error.log; \
                      source /etc/apache2/envvars ; \
                      exec /usr/sbin/apache2 -DFOREGROUND"
stdout_logfile=/var/log/apache2/access.log
stderr_logfile=/var/log/apache2/error.log

autostart=true
autorestart=true
stopasgroup=true
