# assets/etc/supervisor/conf.d/mysql.ini
[program:mysql]
username=mysql
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld

autostart=true
autorestart=true
stopasgroup=true
